# ansible/roles/bot/files/supervisord.conf
[program:bot_app]
directory=/app
command=pipenv run python3 main.py
autostart=true
autorestart=true